# rules/sigma/2026-06-25-macos-gaslight-rust-backdoor.yml
title: macOS.Gaslight Rust Backdoor - LaunchAgent Persistence
id: a1b2c3d4-e5f6-7890-abcd-ef1234567890
status: experimental
description: Detects creation of a LaunchAgent plist with a label masquerading as an Apple system service, consistent with macOS.Gaslight persistence mechanism.
references:
    - https://www.sentinelone.com/labs/macos-gaslight-rust-backdoor-turns-prompt-injection-on-the-analyst-not-the-sandbox/
author: Actioner CTI
date: 2026-06-25
tags:
    - attack.persistence
    - attack.t1543.001
    - attack.t1036.004
logsource:
    category: file_event
    product: macos
detection:
    selection_path:
        TargetFilename|contains: '/Library/LaunchAgents/'
    selection_label:
        TargetFilename|contains: 'com.apple.system.services.activity'
    condition: selection_path and selection_label
falsepositives:
    - Legitimate Apple system services (unlikely to use this exact label)
level: high
---
title: macOS.Gaslight Rust Backdoor - IOPMAssertion Sleep Prevention
id: b2c3d4e5-f6a7-8901-bcde-f12345678901
status: experimental
description: Detects a non-Apple process calling IOPMAssertionCreateWithName to prevent system sleep, a technique used by macOS.Gaslight to maintain persistence.
references:
    - https://www.sentinelone.com/labs/macos-gaslight-rust-backdoor-turns-prompt-injection-on-the-analyst-not-the-sandbox/
author: Actioner CTI
date: 2026-06-25
tags:
    - attack.persistence
    - attack.t1529
logsource:
    category: process_creation
    product: macos
detection:
    selection:
        CommandLine|contains: 'IOPMAssertionCreateWithName'
    filter_apple:
        Image|startswith: '/System/'
    condition: selection and not filter_apple
falsepositives:
    - Legitimate applications that prevent sleep (e.g., video players, download managers)
level: medium
---
title: macOS.Gaslight - Python Stealer Data Collection Archive
id: c3d4e5f6-a7b8-9012-cdef-123456789012
status: experimental
description: Detects creation of the collected_data.zip archive in a temp directory, consistent with the macOS.Gaslight embedded Python stealer staging exfiltration data.
references:
    - https://www.sentinelone.com/labs/macos-gaslight-rust-backdoor-turns-prompt-injection-on-the-analyst-not-the-sandbox/
author: Actioner CTI
date: 2026-06-25
tags:
    - attack.collection
    - attack.t1560.001
    - attack.exfiltration
    - attack.t1041
logsource:
    category: file_event
    product: macos
detection:
    selection:
        TargetFilename|endswith: '/temp/collected_data.zip'
    condition: selection
falsepositives:
    - Legitimate applications creating similarly named archives
level: high
---
title: macOS.Gaslight - Python Runtime Staging from astral-sh
id: d4e5f6a7-b8c9-0123-defa-234567890123
status: experimental
description: Detects download or execution of cpython-3.10 standalone builds from astral-sh/python-build-standalone, used by macOS.Gaslight to stage a Python runtime for its embedded stealer.
references:
    - https://www.sentinelone.com/labs/macos-gaslight-rust-backdoor-turns-prompt-injection-on-the-analyst-not-the-sandbox/
author: Actioner CTI
date: 2026-06-25
tags:
    - attack.execution
    - attack.t1059.006
    - attack.resource_development
    - attack.t1608.001
logsource:
    category: process_creation
    product: macos
detection:
    selection_curl:
        CommandLine|contains|all:
            - 'astral-sh'
            - 'python-build-standalone'
            - 'cpython-3.10'
    condition: selection_curl
falsepositives:
    - Developers legitimately using astral-sh python-build-standalone
level: medium
---
title: macOS.Gaslight - Telegram Bot API C2 Communication
id: e5f6a7b8-c9d0-1234-efab-345678901234
status: experimental
description: Detects process creation with command lines containing Telegram Bot API getUpdates endpoint, consistent with macOS.Gaslight C2 polling mechanism.
references:
    - https://www.sentinelone.com/labs/macos-gaslight-rust-backdoor-turns-prompt-injection-on-the-analyst-not-the-sandbox/
author: Actioner CTI
date: 2026-06-25
tags:
    - attack.command_and_control
    - attack.t1102.002
logsource:
    category: process_creation
    product: macos
detection:
    selection:
        CommandLine|contains|all:
            - 'api.telegram.org'
            - 'getUpdates'
    condition: selection
falsepositives:
    - Legitimate Telegram bots running on macOS
level: medium
