# rules/sigma/2026-06-27-sharkloader-strikeshaker.yml
# Sigma Rules -- SharkLoader / StrikeShark Campaign
# Generated by Actioner CTI -- 2026-06-27
# revision: 2026-06-27-R1 -- applied critic fixes (UUIDs, Sigma 4 rename/downgrade, Sigma 5 comment + confidence)
---
title: SharkLoader DLL Side-Loading via SystemSettings.exe from Non-Standard Path
id: face6b72-9660-46bb-bd91-24db5e6fb8d3
status: experimental
description: Detects execution of SystemSettings.exe from non-standard directories, indicative of SharkLoader DLL side-loading as observed in the StrikeShark campaign.
references:
    - https://securelist.com/strikeshark-campaign/120326/
    - https://thehackernews.com/2026/06/new-sharkloader-malware-deploys-cobalt.html
author: Actioner CTI
date: 2026-06-27
tags:
    - attack.defense_evasion
    - attack.t1574.002
    - attack.execution
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        Image|endswith: '\SystemSettings.exe'
    filter_legitimate:
        Image|startswith:
            - 'C:\Windows\ImmersiveControlPanel\'
            - 'C:\Windows\System32\'
            - 'C:\Windows\SysWOW64\'
    condition: selection and not filter_legitimate
falsepositives:
    - Legitimate software that bundles SystemSettings.exe in non-standard locations
level: high
---
title: SharkLoader Registry Run Key Persistence - MFUpdate
id: 1ec4e74f-ed3a-4b36-8152-1caa1279a61d
status: experimental
description: Detects creation of the MFUpdate registry Run key value pointing to SystemSettings.exe under Identities folder, a persistence mechanism used by SharkLoader in the StrikeShark campaign.
references:
    - https://securelist.com/strikeshark-campaign/120326/
    - https://thehackernews.com/2026/06/new-sharkloader-malware-deploys-cobalt.html
author: Actioner CTI
date: 2026-06-27
tags:
    - attack.persistence
    - attack.t1547.001
logsource:
    category: registry_set
    product: windows
detection:
    selection:
        TargetObject|endswith: '\SOFTWARE\Microsoft\Windows\CurrentVersion\Run\MFUpdate'
        Details|contains: '\Identities\SystemSettings.exe'
    condition: selection
falsepositives:
    - Unlikely
level: critical
---
title: SharkLoader Scheduled Task Creation - Edge Update Masquerade
id: 3e72b4d4-3bc5-4b39-b2d8-35d3e946760c
status: experimental
description: Detects creation of scheduled tasks masquerading as Microsoft Edge updates used by SharkLoader for persistence, executing SystemSettings.exe from suspicious paths.
references:
    - https://securelist.com/strikeshark-campaign/120326/
    - https://thehackernews.com/2026/06/new-sharkloader-malware-deploys-cobalt.html
author: Actioner CTI
date: 2026-06-27
tags:
    - attack.persistence
    - attack.t1053.005
logsource:
    category: process_creation
    product: windows
detection:
    selection_schtasks:
        Image|endswith: '\schtasks.exe'
        CommandLine|contains|all:
            - '/create'
            - 'Edgeupdate'
            - 'SystemSettings.exe'
    condition: selection_schtasks
falsepositives:
    - Unlikely
level: critical
---
title: NTDS Credential Extraction via ntdsutil IFM
id: 4bf1d261-7cbe-4cb9-9845-47c33b076531
status: experimental
description: Detects ntdsutil credential extraction commands using Install From Media (IFM) to dump the Active Directory database. This is a generic credential access technique observed across multiple campaigns including StrikeShark.
references:
    - https://securelist.com/strikeshark-campaign/120326/
    - https://thehackernews.com/2026/06/new-sharkloader-malware-deploys-cobalt.html
author: Actioner CTI
date: 2026-06-27
tags:
    - attack.credential_access
    - attack.t1003.003
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        Image|endswith: '\ntdsutil.exe'
        CommandLine|contains|all:
            - 'ac i ntds'
            - 'ifm'
            - 'create full'
    condition: selection
falsepositives:
    - Legitimate Active Directory backup operations
level: medium
---
title: SharkLoader Dropper Execution - Known Filenames
id: c2f80e9d-4c2f-4eb2-b7bd-edbb55c2c683
status: experimental
description: |
    Detects execution of known SharkLoader dropper filenames disguised as legitimate software installers, including fake Cisco AnyConnect and Google Update binaries.
    Note: The AnyConnect filename intentionally ends in "k9exe" without a dot separator. This matches the actual dropper filename observed in the wild, which omits the file extension dot to evade simple extension-based filters.
references:
    - https://securelist.com/strikeshark-campaign/120326/
    - https://thehackernews.com/2026/06/new-sharkloader-malware-deploys-cobalt.html
author: Actioner CTI
date: 2026-06-27
tags:
    - attack.execution
    - attack.t1204.002
    - attack.defense_evasion
    - attack.t1036.005
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        Image|endswith:
            - '\GoogleUpdateStepup.exe'
            - '\AnyConnect-win-4.10.04071-predeploy-k9exe'
    condition: selection
falsepositives:
    - Legitimate Google Update or Cisco AnyConnect installers with matching names
level: medium
