# rules/sigma/2026-06-07-pypi-hades-cluster-supply-chain.yml
# Hades Cluster / Shai-Hulud PyPI Supply Chain Attack - Sigma Detection Rules
# Generated by Actioner on 2026-06-07
# References:
#   - https://securityonline.info/pypi-supply-chain-attack/
#   - https://www.hendryadrian.com/shai-hulud-descends-to-hades-miasma-worm-campaign-spreads-with-new-pypi-wave/
#   - https://snyk.io/blog/lightning-pypi-compromise-bun-based-credential-stealer/
#   - https://research.jfrog.com/post/shai-hulud-here-we-go-again/
---
title: Suspicious Python .pth Startup Hook File Creation
id: 8e3a1d4f-7b2c-4e9a-b5d6-1f0c8a9e3d7b
status: experimental
description: >
    Detects creation of .pth files in Python site-packages directories, consistent with
    the Hades Cluster / Shai-Hulud supply chain attack that uses *-setup.pth startup hooks
    to execute credential-stealing payloads during Python interpreter initialization.
references:
    - https://securityonline.info/pypi-supply-chain-attack/
    - https://www.hendryadrian.com/shai-hulud-descends-to-hades-miasma-worm-campaign-spreads-with-new-pypi-wave/
    - https://snyk.io/blog/lightning-pypi-compromise-bun-based-credential-stealer/
author: Actioner
date: 2026/06/07
tags:
    - attack.t1546.016
    - attack.t1195.001
logsource:
    category: file_event
    product: linux
detection:
    selection:
        TargetFilename|contains: 'site-packages'
        TargetFilename|endswith: '-setup.pth'
    condition: selection
falsepositives:
    - Legitimate Python packages using .pth files for path configuration
level: high
---
title: Bun Runtime Spawned by Python Process
id: 2f4b8c6a-9d1e-4a3f-b7c5-0e2d6f8a1b4c
status: experimental
description: >
    Detects the Bun JavaScript runtime being spawned as a child of a Python process,
    a key behavior of the Hades Cluster credential stealer that downloads Bun v1.3.13
    to execute an obfuscated JavaScript payload (_index.js / router_runtime.js).
references:
    - https://securityonline.info/pypi-supply-chain-attack/
    - https://www.hendryadrian.com/shai-hulud-descends-to-hades-miasma-worm-campaign-spreads-with-new-pypi-wave/
    - https://snyk.io/blog/lightning-pypi-compromise-bun-based-credential-stealer/
author: Actioner
date: 2026/06/07
tags:
    - attack.t1059.007
    - attack.t1195.001
logsource:
    category: process_creation
    product: linux
detection:
    selection:
        ParentImage|endswith:
            - '/python'
            - '/python3'
            - '/python3.10'
            - '/python3.11'
            - '/python3.12'
            - '/python3.13'
        Image|endswith: '/bun'
        CommandLine|contains:
            - '_index.js'
            - 'router_runtime.js'
    condition: selection
falsepositives:
    - Developers intentionally using Bun from Python build scripts
level: critical
---
title: Hades Cluster Credential Stealer Persistence via Systemd User Service
id: 4c7e9a2b-1d5f-4b8c-a3e6-0f9d2c8b5a1e
status: experimental
description: >
    Detects creation of systemd user services used by the Hades Cluster / Shai-Hulud
    credential stealer for persistence, specifically gh-token-monitor.service and
    pgsql-monitor.service which implement a dead-man switch polling GitHub API.
references:
    - https://securityonline.info/pypi-supply-chain-attack/
    - https://research.jfrog.com/post/shai-hulud-here-we-go-again/
    - https://snyk.io/blog/lightning-pypi-compromise-bun-based-credential-stealer/
author: Actioner
date: 2026/06/07
tags:
    - attack.t1543.002
logsource:
    category: file_event
    product: linux
detection:
    selection:
        TargetFilename|contains: '.config/systemd/user/'
        TargetFilename|endswith:
            - 'gh-token-monitor.service'
            - 'pgsql-monitor.service'
    condition: selection
falsepositives:
    - Unlikely - these service names are specific to this campaign
level: critical
---
title: GitHub API DNS Query from Bun or Python Process
id: 6a8d2e4f-3b1c-4f7a-9d5e-2c0b8a6f1d3e
status: experimental
description: >
    Detects DNS queries to api.github.com from Bun or Python processes.
    The Hades Cluster stealer polls this endpoint for dead-man switch checks
    and exfiltrates credentials via GraphQL mutations, but any developer
    tool using the GitHub API will also trigger this rule.
references:
    - https://securityonline.info/pypi-supply-chain-attack/
    - https://www.hendryadrian.com/shai-hulud-descends-to-hades-miasma-worm-campaign-spreads-with-new-pypi-wave/
    - https://snyk.io/blog/lightning-pypi-compromise-bun-based-credential-stealer/
author: Actioner
date: 2026/06/07
tags:
    - attack.t1567.001
logsource:
    category: dns_query
detection:
    selection_domain:
        QueryName: 'api.github.com'
    selection_process:
        Image|endswith:
            - '/bun'
            - '/python'
            - '/python3'
    condition: selection_domain and selection_process
falsepositives:
    - Legitimate developer tools querying GitHub API
    - CI/CD pipelines interacting with GitHub
    - Any Python or Bun application using GitHub REST or GraphQL API
level: low
